# root/etc/s6-overlay/s6-rc.d/svc-spawn-fcgi/run
#!/usr/bin/with-contenv bash

exec 2>&1 \
    /usr/bin/spawn-fcgi -n -s /run/fcgiwrap.sock -u abc -- /usr/bin/fcgiwrap
