# services/apache2/run
#!/bin/bash

. /etc/apache2/envvars

exec /sbin/setuser abc /usr/sbin/apache2 -k start -DNO_DETACH