# services/apache2/run
#!/bin/bash
exec /usr/sbin/apache2ctl -D FOREGROUND
