# services/smokeping/run
#!/bin/bash

exec /sbin/setuser abc /usr/sbin/smokeping --config="/etc/smokeping/config" --nodaemon
