# services/smokeping/run
#!/bin/bash

exec /usr/sbin/smokeping --config /etc/smokeping/config --nodaemon