# services/smokeping/run
#!/bin/bash

exec "/usr/sbin/smokeping --nodaemon"